# src/lxc/virtlxcd.service.in
[Unit]
Description=Virtualization lxc daemon
Conflicts=libvirtd.service
Requires=virtlxcd.socket
Requires=virtlxcd-ro.socket
Requires=virtlxcd-admin.socket
After=network.target
After=dbus.service
After=apparmor.service
After=local-fs.target
Documentation=man:libvirtd(8)
Documentation=https://libvirt.org

[Service]
Type=notify
ExecStart=@sbindir@/virtlxcd --timeout 120
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure

[Install]
WantedBy=multi-user.target
Also=virtlxcd.socket
Also=virtlxcd-ro.socket
Also=virtlxcd-admin.socket
